# docker/images/www/supervisord/supervisord-memcached.conf
[program:memcached]
command=/usr/bin/memcached -p 11211 -u memcache -m 64 -c 1024
numprocs=1
autostart=true
autorestart=true