# docker/images/www/supervisord/supervisord-nginx.conf
[program:nginx]
command=/scripts/start-nginx.sh
numprocs=1
autostart=true
autorestart=true