# docker/images/www/supervisord/supervisord-php-fpm.conf
[program:php-fpm]
command=/scripts/start-php-fpm.sh
numprocs=1
autostart=true
autorestart=true